# deployment/systemd/atlas-prod.service
# Atlas Production Service
# Location: /etc/systemd/system/atlas-prod.service
#
# Enable: sudo systemctl enable atlas-prod
# Start:  sudo systemctl start atlas-prod
# Logs:   journalctl -u atlas-prod -f

[Unit]
Description=Atlas Infrastructure Platform (Production)
Documentation=https://github.com/DanielTromp/Atlas
After=network-online.target atlas-traefik.service
Wants=network-online.target
Requires=atlas-traefik.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/data/Atlas
ExecStartPre=/usr/bin/podman-compose down --remove-orphans
ExecStart=/usr/bin/podman-compose up
ExecStop=/usr/bin/podman-compose down
Restart=on-failure
RestartSec=10
TimeoutStartSec=180
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-prod

[Install]
WantedBy=multi-user.target
